# frc-tracker.service
[Unit]
Description=FRC Attendance Tracker
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/time-tracker
EnvironmentFile=/home/pi/time-tracker/.env
ExecStart=/home/pi/time-tracker/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
