# tempus.service
[Unit]
Description=Tempus — FRC Attendance Tracker
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/tempus
EnvironmentFile=/home/pi/tempus/.env
ExecStart=/home/pi/tempus/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
